# packer/images/base/files/etc/init/consul_agent_join.conf
start on started consul
stop on stopped consul

task
exec /usr/bin/python /usr/local/bin/consul_join
