# router/files/etc/service/haproxy/run
#!/bin/sh

exec 2>&1
exec /usr/sbin/haproxy-systemd-wrapper -f /etc/haproxy/haproxy.cfg -p /run/haproxy.pid
